# shopify-pipeline.service
[Unit]
Description=AI Shopify Automation Pipeline - autonomous loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ai-shopify-automation-pipeline
ExecStart=/root/ai-shopify-automation-pipeline/loop.sh
Restart=always
RestartSec=30
StandardOutput=append:/root/ai-shopify-automation-pipeline/logs/loop.log
StandardError=append:/root/ai-shopify-automation-pipeline/logs/loop_errors.log

[Install]
WantedBy=multi-user.target
